# agent-platform/desktop/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:vnc]
command=su - agent -c "vncserver :1 -geometry %(ENV_RESOLUTION)s -depth 24 -SecurityTypes VncAuth -fg"
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/vnc.log
stderr_logfile=/var/log/supervisor/vnc-err.log

[program:novnc]
command=websockify --web /usr/share/novnc %(ENV_NOVNC_PORT)s localhost:%(ENV_VNC_PORT)s
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc-err.log

[program:code-server]
command=code-server --bind-addr 0.0.0.0:8443 --auth none /workspace
user=agent
autorestart=true
priority=30
stdout_logfile=/var/log/supervisor/code-server.log
stderr_logfile=/var/log/supervisor/code-server-err.log
